# deploy/systemd/cs2-monitor-web.service
[Unit]
Description=CS2 Monitor Web Server
After=network.target mysql.service
Wants=mysql.service

[Service]
User=cs2monitor
Group=cs2monitor
WorkingDirectory=/opt/cs2-monitor
ExecStart=/opt/cs2-monitor/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
EnvironmentFile=/opt/cs2-monitor/.env

[Install]
WantedBy=multi-user.target
