# deploy/systemd/cs2-monitor-scheduler.service
[Unit]
Description=CS2 Monitor Price Scheduler
After=network.target mysql.service
Wants=mysql.service

[Service]
User=cs2monitor
Group=cs2monitor
WorkingDirectory=/opt/cs2-monitor
Environment=PYTHONPATH=/opt/cs2-monitor
EnvironmentFile=/opt/cs2-monitor/.env
ExecStart=/opt/cs2-monitor/venv/bin/python -m scheduler.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
